# overlay/etc/supervisor/conf.d/bind.conf
[program:bind]
command=/usr/sbin/named -c /etc/bind/named.conf -u bind -f
stdout_events_enabled=true
stderr_events_enabled=true
